# backup.init.conf
# Home Network Backup Service

description     "Home Network Backup Service"
author          "Jo Vandeginste <jo.vandeginste@gmail.com>"

start on (net-device-up
          and local-filesystems
	  and runlevel [2345])
stop on runlevel [016]

respawn
expect daemon

exec /usr/local/backup/backup
